# misc/init.d-ath
#!/sbin/runscript

depend() {
    need mongod
}

start() {
    ebegin "Starting ath"
    sudo -i -u darshan "ath up"
    eend $? "Failed to start ath"
}

stop() {
    ebegin "Stopping ath"
    sudo -i -u darshan "ath down"
    eend $? "Failed to stop ath"
}
